# deploy/heartbit-report.timer
[Unit]
Description=Heartbit daily report timer

[Timer]
OnCalendar=Mon..Fri 09:00
Persistent=true

[Install]
WantedBy=timers.target
